# supervisor.conf
[program:your_app]
command=/usr/local/bin/python3 your_script.py
directory=/app
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/your_app.log
